# config/cron.d/whybug
33  */6   * * *   root /opt/whybug/discovery.js -h all >> /dev/null 2>&1
